# supervisord.conf
[supervisord]
nodaemon=true

[program:wrapdocker]
command=/usr/local/bin/wrapdocker

[program:aspnetdebug]
command=dnx -p project.json aspnet_debug.Server